# bitcoin-bot.service
[Unit]
Description=Bitcoin Price Bot Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/your/bot
ExecStart=/usr/bin/python3 /path/to/your/bot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 